# rpi/etc/cron.d/rpi-ds18b20
# run every 10 minutes
*/10 * * * * root /opt/rpi-ds18b20/logtemp.sh
